# root/etc/services.d/quassel/run
#!/usr/bin/with-contenv bash

chmod -R 777 /quassel

exec \
	s6-setuidgid abc quasselcore \
	-c /quassel -p 4242 --require-ssl
